# ocs-inventory-ng-client.cron
1 12 * * * root cd /usr/sbin && perl ocsinventory-client.pl > /dev/null 2>&1
